# contrib/systemd/aws-dyndns.timer
[Unit]
Description=Sync aws dns with ISP address

[Timer]
OnBootSec=10min
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
